# contrib/deis-backup.timer
[Unit]
Description=deis-backup-timer

[Timer]
OnCalendar=daily

[X-Fleet]
X-ConditionMachineOf=deis-backup.service