# contrib/deis-backup-restore.timer
[Unit]
Description=deis-backup-restore-timer

[Timer]
OnCalendar=daily

[X-Fleet]
X-ConditionMachineOf=deis-backup-restore.service